# deploy/systemd/envsense-web.service.template
[Unit]
Description=EnvSense AI FastAPI Web Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PROJECT_DIR__
Environment="PATH=__VENV_PATH__/bin:/usr/bin:/bin"
ExecStart=__VENV_PATH__/bin/uvicorn fastapi_app:app --host 0.0.0.0 --port 8501
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

